# systemd/iptooleg.timer
[Unit]
Description=ext ip to oleg

[Timer]
OnBootSec=30
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target

