# systemd/iptooleg.service
[Unit]
Description=ext ip to oleg

[Service]
ExecStart=/opt/sendip.sh
