# systemd/ws.service
[Service]
ExecStart=/usr/bin/node server.js
WorkingDirectory=/usr/share/nginx/html/
User=http
Restart=always

[Install]
WantedBy=multi-user.target
